# docker/crontab
# 每4小时抓取一次新闻
0 */4 * * * cd /app/scripts && /app/venv/bin/python news_scraper.py >> /app/logs/scraper.log 2>&1

# 每天凌晨3点清理旧日志（保留7天）
0 3 * * * find /app/logs -name "*.log" -type f -mtime +7 -delete

# 每周日凌晨4点备份数据库
0 4 * * 0 cp /app/public/news.db /app/public/news.db.bak.$(date +\%Y\%m\%d) && find /app/public -name "news.db.bak.*" -type f -mtime +30 -delete

# 保持 crontab 运行的空行